# deploy/systemd/gscale-bot.service
[Unit]
Description=GScale Zebra telegram bot worker
After=network-online.target gscale-scale.service
Wants=network-online.target
Requires=gscale-scale.service

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__PREFIX__
EnvironmentFile=-__PREFIX__/config/bot.env
ExecStart=__PREFIX__/bin/bot
Restart=always
RestartSec=1
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
