# deploy/systemd/gscale-scale.service
[Unit]
Description=GScale Zebra scale worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__PREFIX__
EnvironmentFile=-__PREFIX__/config/scale.env
ExecStart=__PREFIX__/bin/scale --no-bot --no-bridge --device $SCALE_DEVICE --zebra-device $ZEBRA_DEVICE --bridge-state-file $BRIDGE_STATE_FILE
Restart=always
RestartSec=1
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
